# crontab
*/5 * * * * root python /root/ldap-ctd.tn.corp.py
